# _linux/shell/crontab.cron
00 02 * * * /bin/bash /my/backup/mysql-auto-backup.sh
